# systemd/mull-wg-watcher.path
[Path]
PathModified=/var/tmp/mull-wg/loc

[Install]
WantedBy=multi-user.target
WantedBy=network-online.target