# systemd/mull-wg-serv.service
[Unit]
Description=Mullvad server list update service

[Service]
Type=oneshot

ExecStart=python3 /var/mull-wg/scripts/fetch_servers.py
