# systemd/mull-wg-serv.timer
[Unit]
Description=Mullvad server list update service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target